# .configs/etc/systemd/system/nuxt-watcher.path
[Path]
PathModified=/var/www/somdomato.com

[Install]
WantedBy=multi-user.target