# .configs/etc/systemd/system/nuxt-watcher.service
[Unit]
Description=Som do Mato API Server Watcher

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart nuxt

[Install]
WantedBy=multi-user.target